# templates/usr/lib/systemd/system/riemann-dash.service.erb
[Unit]
Description=Riemann Dashboard
Wants=basic.target
After=basic.target network.target

[Service]
ExecStart=<%= @riemann::dash_path %> /etc/riemann-dash.rb
KillMode=process

[Install]
WantedBy=multi-user.target
